# etc/profile.d/ti-c7x-cgt.sh
export C7X_C_DIR="/usr/share/ti/cgt-c7x/include;/usr/share/ti/cgt-c7x/lib"
